# deploy/systemd/jsnetworkcorp-github-actions-ssh-allowlist.service
[Unit]
Description=Refresh GitHub Actions SSH allowlist for JSNetworkCorp portfolio deploy
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update-github-actions-ssh-allowlist
